# deb/DEBIAN/postinst
#!/bin/bash
set -e

systemctl daemon-reload
systemctl enable cloudreflect.service