# deb/DEBIAN/prerm
#!/bin/bash
set -e

systemctl disable cloudreflect.service
systemctl stop cloudreflect.service